# config/includes.chroot_after_packages/etc/skel/.config/autostart/firstlogin.desktop
[Desktop Entry]
Version=1.0
Type=Application
Name=first login setup
Comment=
Exec=first-login-setup
Icon=
Path=
Terminal=false
StartupNotify=false
